# systemd/gopher.socket
[Unit]
Description=Gopher server socket

[Socket]
ListenStream=70
Accept=true

[Install]
WantedBy=sockets.target

